# cot-monitor.service
# cot-monitor.service
# Copy to /etc/systemd/system/cot-monitor.service
# Then: sudo systemctl daemon-reload && sudo systemctl enable cot-monitor.timer

[Unit]
Description=COT Monitor — Weekly CFTC Report via Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/cot-monitor
ExecStart=/home/cot-monitor/venv/bin/python cot_runner.py
StandardOutput=journal
StandardError=journal

# The script polls for CFTC data for up to 4 hours before giving up.
# This timeout must exceed MAX_WAIT_SEC (14400s) + pipeline runtime (~5 min).
TimeoutStartSec=18000

[Install]
WantedBy=multi-user.target
