# cot-monitor.service
# cot-monitor.service
# Copy to /etc/systemd/system/cot-monitor.service
# Then: sudo systemctl daemon-reload && sudo systemctl enable cot-monitor.timer

[Unit]
Description=COT Monitor — Weekly CFTC Report via Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=YOUR_VPS_USER
WorkingDirectory=/home/YOUR_VPS_USER/cot-monitor
ExecStart=/home/YOUR_VPS_USER/cot-monitor/venv/bin/python cot_runner.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
